# roles/app/templates/sidekiq.service.j2
[Unit]
After=network.target
Description={{ app_sidekiq_service_description }}
Documentation=http://github.com/ngonzalez

[Service]
Type=simple
WorkingDirectory={{ app_path_backend }}
ExecStart=/usr/bin/bundle3.0 exec sidekiq -e production
EnvironmentFile={{ app_path_backend }}/.env
PIDFile={{ app_path_backend }}{{ app_sidekiq_pid_path }}
Restart=always
User={{ ansible_user }}
Group={{ ansible_user }}
SyslogIdentifier=sidekiq
LimitNOFILE=65535
TimeoutStopSec=0

[Install]
WantedBy=multi-user.target
Alias={{ app_sidekiq_service_name }}.service
